# microservice/service/navigation/debian/data/etc/systemd/system/navigation.service
; Navigation systemd service unit file
; Place in e.g. `/etc/systemd/system/navigation.service`, then `systemctl daemon-reload` and `service navigation start`.

[Unit]
Description=Navigation
Requires=network.target
After=network.target

[Service]
ExecStart=/bin/sh -c '/usr/share/shop/navigation/bin/navigation > /var/log/navigation.log 2>&1'
EnvironmentFile=/etc/default/navigation
Type=simple
WorkingDirectory=/usr/share/shop/navigation/bin/
User=root
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
