# AppStash/microservice/service/product/debian/control/postinst
#!/bin/sh
set -e

chmod +x /usr/share/shop/product/bin/product
chmod +x /etc/init.d/product

update-rc.d product defaults

service product force-reload